# deployment/crm-backend.socket
[Unit]
Description=CRM Backend Gunicorn Socket

[Socket]
ListenStream=/run/crm-backend.sock

[Install]
WantedBy=sockets.target
